# catkin_ws/devel/etc/catkin/profile.d/50.segway_description.sh
/home/chels/catkin_ws/devel/.private/segway_description/etc/catkin/profile.d/50.segway_description.sh