# catkin_ws/devel/etc/catkin/profile.d/bwi_planning_common.sh
/home/chels/catkin_ws/devel/.private/bwi_planning_common/etc/catkin/profile.d/bwi_planning_common.sh