# catkin_ws/devel/etc/catkin/profile.d/bwi_knowledge_representation.sh
/home/chels/catkin_ws/devel/.private/bwi_knowledge_representation/etc/catkin/profile.d/bwi_knowledge_representation.sh