# catkin_ws/devel/etc/catkin/profile.d/knowledge_representation.zsh
/home/chels/catkin_ws/devel/.private/knowledge_representation/etc/catkin/profile.d/knowledge_representation.zsh