# catkin_ws/devel/.private/segway_v3_config/etc/catkin/profile.d/50.segway_config.bash
. "/home/chels/catkin_ws/src/segway_v3/segway_v3_config/segway_config.bash"
